# ad-reporting.service
[Unit]
Description=Ad Reporting AdonisJS Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/ad-reporting
ExecStart=/usr/bin/node bin/server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3333

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ad-reporting

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/ad-reporting

[Install]
WantedBy=multi-user.target
